# templates/all/etc/systemd/system/tailscaled.service.d/override.conf
[Service]
Environment=TS_NO_NETWORK_MANAGER=true
Environment=TS_DEBUG_NETSTACK_IDLE_TIMEOUT=1s
