# templates/all/etc/systemd/system/snmpd.service.d/local.conf.j2
[Service]
ExecStart=
ExecStart=/usr/sbin/snmpd -Lsd -f -u root -g root -I -smux,mteTrigger,mteTriggerConf -p /run/snmpd.pid
NotifyAccess=all
